# retralabs-webhook.service
[Unit]
Description=RetraLabs Deploy Webhook
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/retralabs
EnvironmentFile=/opt/retralabs/.env
Environment=APP_DIR=/opt/retralabs
Environment=WEBHOOK_PORT=9000
Environment=WEBHOOK_HOST=127.0.0.1
ExecStart=/usr/bin/node /opt/retralabs/deploy-webhook.js
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target


